# armbian_config/userpatches/overlay/etc/systemd/system/serial-getty@ttyS0.service.d/override.conf
# Put term in dumb mode (needed to talk serial with Flipper zero)

[Service]
Environment=TERM=dumb
